# files/main/usr/lib/systemd/system/flatpak-essential-install.service
[Unit]
Description=Install Essential Flatpaks
# Run after network is online but before flatpak-preinstall
After=network-online.target
Before=flatpak-preinstall.service graphical.target
Wants=network-online.target
ConditionPathExists=/usr/bin/flatpak
ConditionPathExists=!/var/lib/flatpak/essential-installed
Documentation=https://github.com/alyraffauf/kyanite
# Limit retries to avoid infinite loops - matches flatpak-preinstall
StartLimitIntervalSec=600
StartLimitBurst=3

[Service]
Type=oneshot
# Install essential Flatpaks that are critical for offline use
ExecStartPre=/usr/bin/flatpak remote-add --if-not-exists --system flathub https://flathub.org/repo/flathub.flatpakrepo
ExecStart=/usr/bin/flatpak install --system --noninteractive -y flathub org.mozilla.firefox
ExecStart=/usr/bin/flatpak install --system --noninteractive -y flathub io.github.kolunmi.Bazaar
# Mark as complete so this only runs once
ExecStartPost=/usr/bin/touch /var/lib/flatpak/essential-installed
RemainAfterExit=yes
# Retry on failure (e.g., network issues)
Restart=on-failure
RestartSec=60
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
